# scripts/codex-setup.service
# codex-setup.service
[Unit]
Description=Validate environment using docs/setup.md before networking
DefaultDependencies=no
Before=network-pre.target
Wants=network-pre.target

[Service]
Type=oneshot
ExecStart=/usr/local/bin/codex -q -a full-auto "doctor docs/setup.md"
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
